# searxng.service
[Unit]
Description=SearXNG self-hosted search (MCP backend)
Documentation=https://github.com/searxng/searxng-docker
After=docker.service network-online.target
Requires=docker.service
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=%h/Github/searxng-tool

ExecStart=/usr/bin/docker compose up -d --wait
ExecStop=/usr/bin/docker compose down --timeout 20

# Kill only the direct docker compose process, not the entire cgroup.
# This prevents systemd from racing with docker's own container stop logic.
KillMode=process
KillSignal=SIGTERM
FinalKillSignal=SIGKILL

# Do NOT restart during shutdown — oneshot + Restart=on-failure conflicts with
# the stop sequence and can cause the service to respawn while systemd is trying
# to shut it down, causing the visible shutdown hang.
# Use the systemd service start/stop commands manually if restart is needed.

TimeoutStartSec=120
# Must be > docker compose --timeout (20s) + overhead. Set to 45s.
TimeoutStopSec=45

[Install]
WantedBy=default.target
